# python/shutdown/shutdown.service

[Unit]
Description=shutdown

[Service]
ExecStart=/home/pi/work/raspi/python/shutdown/shutdown.py

[Install]
WantedBy=multi-user.target
